# raspifan.service.template
# This file was generated automatically by raspifan.

[Unit]
Description = Raspberry Pi fan controller

ConditionPathExists = {{ .TempPath }}

[Service]
Type = notify-reload
ExecStart = {{ .SelfPath }}

[Install]
WantedBy = multi-user.target

# vim: set ft=systemd:
